# systemd/trade-monitor.service
[Unit]
Description=Agentic Crypto Trade Monitor (Protective Stops)
After=trade-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trade_server
Environment="PATH=/opt/trade_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/opt/trade_server/.env
ExecStart=/opt/trade_env/bin/python services/monitor_service.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
